# server/triggeriq.conf
[program:web-1]
command=/bin/sh -c 'gunicorn routes:app -b 127.0.0.1:8000 --reload'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/web-1.log
stderr_logfile=/var/log/triggeriq/web-1.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5000

[program:change_feed-1]
command=/bin/sh -c 'python -u change_feed.py'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/change_feed-1.log
stderr_logfile=/var/log/triggeriq/change_feed-1.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5100

[program:worker-1]
command=/bin/sh -c 'python -u worker.py -c 5'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/worker-1.log
stderr_logfile=/var/log/triggeriq/worker-1.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5200

[program:worker-2]
command=/bin/sh -c 'python -u worker.py -c 5'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/worker-2.log
stderr_logfile=/var/log/triggeriq/worker-2.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5201

[program:worker-3]
command=/bin/sh -c 'python -u worker.py -c 5'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/worker-3.log
stderr_logfile=/var/log/triggeriq/worker-3.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5202

[program:worker-4]
command=/bin/sh -c 'python -u worker.py -c 5'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/worker-4.log
stderr_logfile=/var/log/triggeriq/worker-4.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5203

[program:worker-5]
command=/bin/sh -c 'python -u worker.py -c 5'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/worker-5.log
stderr_logfile=/var/log/triggeriq/worker-5.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5204

[program:worker-6]
command=/bin/sh -c 'python -u worker.py -c 5'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/worker-6.log
stderr_logfile=/var/log/triggeriq/worker-6.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5205

[program:worker-7]
command=/bin/sh -c 'python -u worker.py -c 5'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/worker-7.log
stderr_logfile=/var/log/triggeriq/worker-7.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5206

[program:worker-8]
command=/bin/sh -c 'python -u worker.py -c 5'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/worker-8.log
stderr_logfile=/var/log/triggeriq/worker-8.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5207

[program:worker-9]
command=/bin/sh -c 'python -u worker.py -c 5'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/worker-9.log
stderr_logfile=/var/log/triggeriq/worker-9.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5208

[program:worker-10]
command=/bin/sh -c 'python -u worker.py -c 5'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/worker-10.log
stderr_logfile=/var/log/triggeriq/worker-10.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5209

[program:clock-1]
command=/bin/sh -c 'python -u clock.py '
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/clock-1.log
stderr_logfile=/var/log/triggeriq/clock-1.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5300

[program:brunch-1]
command=/bin/sh -c 'sh -c 'cd ../client && brunch watch --server''
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/brunch-1.log
stderr_logfile=/var/log/triggeriq/brunch-1.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5400

[program:tests-1]
command=/bin/sh -c 'nosetests -v --with-watch --nocapture'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/tests-1.log
stderr_logfile=/var/log/triggeriq/tests-1.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5500

[program:splash-1]
command=/bin/sh -c 'python -m splash.server --port=8950'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/splash-1.log
stderr_logfile=/var/log/triggeriq/splash-1.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5600

[program:rq_dashboard-1]
command=/bin/sh -c 'rq-dashboard'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/triggeriq/rq_dashboard-1.log
stderr_logfile=/var/log/triggeriq/rq_dashboard-1.error.log
user=root
directory=/data/triggeriq/current/server
environment=PORT=5700

[group:triggeriq]
programs=web-1,change_feed-1,worker-1,worker-2,worker-3,worker-4,worker-5,worker-6,worker-7,worker-8,worker-9,worker-10,clock-1,brunch-1,tests-1,splash-1,rq_dashboard-1
